# roles/filesystem/files/btrfs-scrub@.timer
[Unit]
Description=Monthly btrfs scrub on %f

[Timer]
OnCalendar=monthly
RandomizedDelaySec=1 week
Persistent=true

[Install]
WantedBy=timers.target
